# karte/crontab
*/5 8-21 * * *  /home/pi/metar/karte/refresh.sh
*/5 22 * * *     /home/pi/metar/karte/lightsoff.sh
